# tests/examples/RA/userver/bin/RA
#! /bin/sh                                         
#                                                  
### BEGIN INIT INFO                                
# Provides: userver_tcp
# Description: Start Marvelous UServer TCP         
### END INIT INFO                                  

#UTRACE="0 10M 0"
#UOBJDUMP="0 100k 10"
#USIMERR="error.accept"
 export UTRACE UOBJDUMP USIMERR

# Check for existence of needed config file and read it

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
LD_LIBRARY_PATH=/usr/lib:/usr/local/lib

exe=/usr/local/bin/userver_tcp
confdir=/srv/userver/etc

test -x $exe || exit 0
test -d $confdir || exit 0

conffiles=$confdir/RA.cfg

case "$1" in    
    start)
        for conffile in $conffiles ; do
                echo -n "Starting Marvelous UServer TCP ($conffile) "
                                                   $exe -c $conffile &> /dev/null &
                sleep 1 ; kill -0 $! &> /dev/null
        done
        ;;
    stop)
        for conffile in $conffiles ; do
                echo -n "Stopping Marvelous UServer TCP ($conffile) "
                pkill -f   "$exe -c $conffile" &> /dev/null
        done
        ;;
    status)
        for conffile in $conffiles ; do
                echo -n "Checking for Marvelous UServer TCP ($conffile) "
                #pgrep -f   "$exe -c $conffile" &> /dev/null
                pgrep -f   "$exe -c $conffile"
        done
        ;;
    kill)
        echo -n "Killing Marvelous UServer TCP"
        pkill -f   "$exe"
        ;;
    try-restart)
        $0 stop && $0 start
        ;;
    restart)
        $0 stop
        $0 start
        ;;
    *)
        echo "Usage: $0 {start|stop|status|try-restart|restart}"
        exit 1
esac
