# tests/examples/RA/userver/bin/RA.6
#! /bin/sh
### BEGIN INIT INFO
# Provides:          RA
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: RA station
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/srv/userver/bin/userver_tcp
DAEMON_ARGS="-c /srv/userver/etc/RA.cfg"
NAME=RA
PIDFILE=/var/run/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME
DESC="RA station"

#UTRACE="0 10M 0"
#UOBJDUMP="0 100k 10"
#USIMERR="error.accept"
 export UTRACE UOBJDUMP USIMERR

unset TMPDIR

test -x $DAEMON || exit 0

mkdir -p /var/run /var/log /tmp/scan /tmp/mail

if [ -r /srv/RA/etc/environment.conf ]; then
  . /srv/RA/etc/environment.conf
fi

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
. /lib/lsb/init-functions

case "$1" in
  start)

	# Return
   #   0 if daemon has been started
   #   1 if daemon was already running
   #   2 if daemon could not be started

	log_begin_msg "Starting $DESC: $NAME"

   	start-stop-daemon --start --quiet --pidfile $PIDFILE --background --exec $DAEMON -- \
          $DAEMON_ARGS \
          || return 2

	log_end_msg $?

   # Add code here, if necessary, that waits for the process to be ready
   # to handle requests from services started subsequently which depend
   # on this one.  As a last resort, sleep for some time.

	;;
  stop)
        # Return
        #   0 if daemon has been stopped
        #   1 if daemon was already stopped
        #   2 if daemon could not be stopped
        #   other if a failure occurred

	log_begin_msg "Stopping $DESC: $NAME"

        start-stop-daemon --stop --pidfile $PIDFILE --name userver_tcp

        RETVAL="$?"

  	log_end_msg $RETVAL

        [ "$RETVAL" != 0 ] && return $RETVAL 

        # Many daemons don't delete their pidfiles when they exit.
        rm -f $PIDFILE
        return 0 
	;;

  reload|force-reload)
       log_begin_msg "Reloading $DESC: $NAME"
       start-stop-daemon --stop --quiet --pidfile $PIDFILE --name userver_tcp --signal 1 && success=1
       log_end_msg $?
       ;;
  restart)
	log_begin_msg "Restarting $DESC: $NAME"
	if start-stop-daemon --stop --quiet --retry 5 --oknodo --pidfile $PIDFILE --name userver_tcp; then
		start-stop-daemon --start --quiet --pidfile $PIDFILE --background --exec $DAEMON -- $DAEMON_ARGS
	fi
	log_end_msg $?
	;;
  status)
	echo -n "Status of $DESC: "
	if [ ! -r "$PIDFILE" ]; then
		echo "$NAME is not running."
		exit 3
	fi
	read pid < "$PIDFILE"
	if ps -p "$pid" > /dev/null 2>&1; then
		echo "$NAME is running."
		exit 0
	else
		echo "$NAME is not running but $PIDFILE exists."
		exit 1
	fi
	;;
  *)
	N=/etc/init.d/${0##*/}
	echo "Usage: $N {start|stop|restart|force-reload|status}" >&2
	exit 1
	;;
esac

exit 0
